# debian/prerm
#!/bin/sh
# prerm script for nordvpn-sysvinit-compat
#
set -e
# stop daemon
if [ -x "/etc/init.d/nordvpns" ] && [ "$1" = remove ]; then
	invoke-rc.d nordvpns stop || exit $?
fi

# End 

exit 0
